# systemd/app.target
[Unit]
Description=Example Target Unit

[Install]
WantedBy=multi-user.target
